# deploy/echochat-frontend-vps.service
[Unit]
Description=EchoChat Frontend (Next.js)
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/echochat/frontend
ExecStart=/usr/bin/node /opt/echochat/frontend/.next/standalone/server.js
Restart=always
RestartSec=5

# Node Environment
Environment=NODE_ENV=production
Environment=PORT=3001

# API URL
Environment=NEXT_PUBLIC_API_URL=http://localhost:8001

# NextAuth Configuration
# IMPORTANT: Replace these values with your actual credentials!
Environment=NEXTAUTH_SECRET=your-secret-here-replace-me
Environment=NEXTAUTH_URL=https://echochat.iafluence.cloud

# Google OAuth Configuration
# IMPORTANT: Replace these values with your actual Google OAuth credentials!
Environment=GOOGLE_CLIENT_ID=your-client-id.apps.googleusercontent.com
Environment=GOOGLE_CLIENT_SECRET=your-client-secret-replace-me

# Admin Email Whitelist
Environment=ALLOWED_ADMIN_EMAILS=suan.tay.job@gmail.com,j.wallut.pro@gmail.com

[Install]
WantedBy=multi-user.target
